# quadlet/caddy.container.in
[Unit]
Description=Caddy HTTPS reverse proxy
Wants=network-online.target llm-net-network.service
After=network-online.target llm-net-network.service caddy-build.service llmster.service

[Container]
ContainerName=caddy
Image=caddy.build
Network=llm-net.network
Volume=caddy-data.volume:/data
Volume=caddy-config.volume:/config
PublishPort=1243:443

EnvironmentFile=@SYSCONFDIR@/llmster/llmster.env

PodmanArgs=--log-driver=json-file --log-opt=max-size=10m --log-opt=max-file=3

[Service]
Restart=on-failure

[Install]
WantedBy=multi-user.target default.target
