# quadlet/llmster.container.in
[Unit]
Description=LM Studio headless server (llmster)
Wants=network-online.target llm-net-network.service
After=network-online.target llm-net-network.service llmster-build.service

[Container]
ContainerName=llmster
Image=llmster.build
Network=llm-net.network
Volume=lmstudio-data.volume:/root/.lmstudio
PublishPort=1234:1234

AddDevice=/dev/dri
AddDevice=/dev/kfd
AddDevice=/dev/fuse

EnvironmentFile=@SYSCONFDIR@/llmster/llmster.env
Environment=LIBVA_DRIVER_NAME=radeonsi
Environment=LLMSTER_PORT=1234
Environment=JIT_LOADING=true
Environment=AUTO_UNLOAD=true
Environment=UNLOAD_IDLE_TIME=300
Environment=OLLAMA_ORIGINS=*
Environment=LMS_SERVER_HOST=0.0.0.0

HealthCmd=/healthcheck.sh
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

PodmanArgs=--memory-reservation=8g --memory=24g --cpus=4 --log-driver=json-file --log-opt=max-size=10m --log-opt=max-file=3

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target
